# worker/worker/templates/fieldnotes-backup.service
[Unit]
Description=fieldnotes daily backup

[Service]
Type=oneshot
ExecStart={{EXEC_START}}
